# docker/worker/crontab
* * * * * php /var/www/html/artisan schedule:run 2>&1
